# lotus-exporter-farcaster/lotus-exporter-farcaster.service
[Unit]
Description=Lotus Exporter Farcaster

[Service]
ExecStart=/usr/local/bin/lotus-exporter-farcaster.py
User=lotus