# lotus-exporter-farcaster/lotus-exporter-farcaster.socket
[Unit]
Description=Lotus Exporter Farcaster

[Socket]
ListenStream=9101

[Install]
WantedBy=default.target