# lotus-exporter-farcaster/lotus-exporter-farcaster.cron
# EXPORT LOTUS TO NODE-EXPORTER / PROMETHEUS
* * * * *   LOTUS_USER if TARGET_BIN_PATH/lotus-exporter-farcaster.py CUSTOM_CONFIG_ARG > /var/lib/prometheus/node-exporter/farcaster.prom.$$;  then mv /var/lib/prometheus/node-exporter/farcaster.prom.$$ /var/lib/prometheus/node-exporter/farcaster.prom; else rm /var/lib/prometheus/node-exporter/farcaster.prom.$$; fi 

