# lotus-exporter-farcaster/lotus-exporter-farcaster.cron
* * * * *   lotus if lotus-farcaster/lotus-exporter-farcaster/lotus-exporter-farcaster.py > /var/lib/prometheus/node-exporter/farcaster.prom.$$;  then mv /var/lib/prometheus/node-exporter/farcaster.prom.$$ /var/lib/prometheus/node-exporter/farcaster.prom; else rm /var/lib/prometheus/node-exporter/farcaster.prom.$$; fi 

